# systemd/vending.service
[Unit]
Description=Control the vending machine
After=network.target

[Service]
Environment="DB_HOST=" "DB_USER=" "DB_PASSWORD=" "CARD_AUTH_USER=" "CARD_AUTH_PASS="
WorkingDirectory=/opt/vending
ExecStart=/usr/bin/node bin/vend.js
Type=notify
NotifyAccess=all
TimeoutSec=5
WatchdogSec=10
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
